# config/systemd/user/organize-files.service
[Unit]
Description=Run the file organizer script

[Service]
Type=oneshot
ExecStart=%h/.local/share/nmde/bin/organize_files
